# systemd/mediatoad.service
[Unit]
Description=MediaToad
After=network.target

[Service]
Environment="PATH=/usr/bin:/bin"
Type=simple
User=librarian
ExecStart=/usr/bin/java \
          -XX:+PerfDisableSharedMem \
          -XX:-UsePerfData \
          -jar /home/librarian/mediatoad/mediatoad.jar \
          --idfile /home/librarian/mediatoad/systemid \
          --tree /home/librarian/mediatoad/tree \
          --db /home/librarian/mediatoad/db \
          --thumbs /home/librarian/mediatoad/thumbs
Restart=always
RestartSec=60
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
